# experiment_scripts/geni/setup_scripts/nco.service
[Unit]
Description=NCO systemd service.

[Service]
Type=simple
ExecStart='/users/dflukasz/software_defined_customization/NCO/NCO.py'

[Install]
WantedBy=multi-user.target
